# examples/adn-monitor.service
# ADN Systems Monitor – systemd unit example
# Copy to /etc/systemd/system/adn-monitor.service and adjust paths.
# Enable: sudo systemctl enable adn-monitor && sudo systemctl start adn-monitor

[Unit]
Description=ADN Systems Monitor

[Service]
User=root
Type=simple
Restart=always
RestartSec=3
SyslogIdentifier=adn-monitor
WorkingDirectory=/opt/adn-monitor/monitor
ExecStart=/usr/bin/python3 /opt/adn-monitor/monitor/monitor.py

[Install]
WantedBy=multi-user.target
